# udev/rules.d/99-monitor-hotplug.rules
ACTION=="change", SUBSYSTEM=="drm", ENV{HOTPLUG}=="1", RUN+="/home/bcherukuri/dotfiles/udev/scripts/i3/reload.sh"
